# system/mboard-update.service
[Unit]
Description=Missionary Board Update Service

[Service]
Type=oneshot
ExecStart=%h/missionary-board/system/mboard-update.sh
