# system/mboard.service
[Unit]
Description=Missionary Board Server

[Service]
Restart=always
WorkingDirectory=/home/orangepi/missionary-board
ExecStart=/home/$username/.local/bin/uv run uvicorn --app-dir ./src --log-config ./conf/uvicorn.logger.json mboard.main:app

[Install]
WantedBy=multi-user.target
